# config/ubuntu-12.04/cron.daily/newdjangosite-daily-database
#!/bin/bash

echo 'VACUUM ANALYZE;' | psql newdjangosite-daily
/var/www/python/newdjangosite-daily/venv/bin/python /var/www/python/newdjangosite-daily/web/manage_daily.py clearsessions

