# config/ubuntu-12.04/cron.daily/newdjangosite-staging-database
#!/bin/bash

echo 'VACUUM ANALYZE;' | psql newdjangosite_staging
/var/www/python/newdjangosite-staging/venv/bin/python /var/www/python/newdjangosite-staging/web/manage_staging.py clearsessions

