# config/ubuntu-16.04/cron.daily/newdjangosite-dev-database
#!/bin/bash

echo 'VACUUM ANALYZE;' | psql newdjangosite_dev
/var/www/python/newdjangosite-dev/venv/bin/python /var/www/python/newdjangosite-dev/web/manage_dev.py clearsessions

